# supervisord.conf
[supervisord]
nodaemon=true

[program:mlflow]
command=/app/mlflow.sh

[program:nginx]
command=/app/webserver.sh
